# target/linux/imx/base-files/etc/uci-defaults/50_disable-removing-wdev
. /lib/functions.sh

case "$(board_name)" in
	aqara,zhwg11lm |\
	xiaomi,dgnwg05lm )
		sed -i 's/iw dev "$wdev" del/# iw dev "$wdev" del/' /lib/netifd/wireless/mac80211.sh
		sed -i 's/nl80211.request(nl80211.const.NL80211_CMD_DEL_INTERFACE/\/\/ nl80211.request(nl80211.const.NL80211_CMD_DEL_INTERFACE/' /usr/share/hostap/common.uc
	;;
esac

exit 0
